# system/services/xnav-firstboot.service
[Unit]
Description=XNav First-Boot Package Installation
Documentation=https://github.com/realaaravdas/Limelight3-XNav
After=network-online.target xnav-expand-rootfs.service
Wants=network-online.target
# Only run if the Python venv has not been created yet.
# On every subsequent boot this condition fails and systemd skips the unit
# instantly, so xnav-vision and xnav-dashboard start without any delay.
ConditionPathExists=!/opt/xnav/venv
Before=xnav-vision.service xnav-dashboard.service

[Service]
Type=oneshot
RemainAfterExit=yes
ExecStart=/bin/bash /etc/xnav/first_boot.sh
TimeoutStartSec=600
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target
